# system/services/bt-pairing.service
[Unit]
Description=Bluetooth Auto-Pairing Agent
After=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
ExecStartPre=/usr/bin/bluetoothctl power on
ExecStartPre=/usr/bin/bluetoothctl discoverable on
ExecStartPre=/usr/bin/bluetoothctl pairable on
ExecStartPre=/bin/sleep 1
ExecStart=/usr/bin/bt-agent -c NoInputNoOutput
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
