# system/services/bt_gps_receiver.service
[Unit]
Description=PwnHub Bluetooth GPS Receiver
After=bluetooth.target
Requires=bluetooth.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/bin/bt_gps_receiver
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal

# Bluetooth compatibility
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/dbus/system_bus_socket

[Install]
WantedBy=multi-user.target
