# system/services/pwnaui.service
[Unit]
Description=PwnaUI - High-Performance UI Renderer for Pwnagotchi
After=local-fs.target
Before=pwnagotchi.service

[Service]
Type=simple
ExecStart=/usr/local/bin/pwnaui -D waveshare2in13_v4 -p -B -w
ExecStartPost=/usr/local/bin/pwnaui-set-theme
Restart=on-failure
RestartSec=5
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=false
PrivateTmp=false
ReadWritePaths=/var/run /home/pi /tmp /var/lib/pwnagotchi
MemoryMax=32M
CPUQuota=20%

[Install]
WantedBy=multi-user.target
