# pwnaui_BACKUP_20260130_123704/pwnaui.service
[Unit]
Description=PwnaUI - High-Performance UI Renderer for Pwnagotchi
Documentation=https://github.com/pwnagotchi/pwnaui
After=local-fs.target
Before=pwnagotchi.service

[Service]
Type=simple
ExecStart=/usr/local/bin/pwnaui -D waveshare2in13_v4 -p -b
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=false
PrivateTmp=true
ReadWritePaths=/var/run

# Resource limits
MemoryMax=32M
CPUQuota=20%

[Install]
WantedBy=multi-user.target
