# pwnaui_BACKUP_20260130_123704/nexmon_fixes/systemd/nexmon-watchdog.service
[Unit]
Description=Nexmon WiFi Watchdog and Recovery Service
Documentation=https://github.com/seemoo-lab/nexmon/issues/335
After=network.target pwnagotchi.service
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/wifi_recovery.sh watch
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexmon-watchdog

# Environment
Environment="NEXMON_INTERFACE=wlan0mon"
Environment="NEXMON_BASE_INTERFACE=wlan0"

# Resource limits
Nice=-5
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Security hardening
NoNewPrivileges=no
ProtectSystem=false
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
